# speedportpro/rootfs/etc/services.d/speedportpro/run
#!/usr/bin/with-contenv bashio
# Dein Start-Skript
exec /run.sh
